# service_v2/energy-grid.service
[Unit]
Description=Energy Grid Simulation Service
After=network.target

[Service]
User=grid
Group=grid
Type=simple
WorkingDirectory=/opt/energy-grid/
ExecStart=/opt/energy-grid/energy-grid
Restart=no
RestartSec=5
StandardOutput=file:/var/log/energy-grid/grid.log
StandardError=file:/var/log/energy-grid/grid.log

[Install]
WantedBy=multi-user.target
